# maintenance-mode.service
[Service]
ExecStart=/opt/maintenance-mode/maintenance-mode.sh
Restart=on-failure

[Install]
WantedBy=default.target
